# systemd/python_crawler.timer
Description=Run foo weekly and on boot

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target